# github-issues-bot.service
[Unit]
Description=GitHub Issues Bot
After=network.target

[Service]
Type=simple
User=github-bot
WorkingDirectory=/opt/github-issues-bot
ExecStart=/bin/bash -c 'source /home/github-bot/.nvm/nvm.sh && nvm use 18 && node dist/index.js'
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=github-issues-bot
Environment=NODE_ENV=production
Environment=NVM_DIR=/home/github-bot/.nvm

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/github-issues-bot/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

[Install]
WantedBy=multi-user.target
